# musicbot.service
[Unit]
After=network.service
[Service]
ExecStart=/home/ubuntu/MusicBot/run.py
